# systemd/split-ssh-agent@.service
[Unit]
Description=split-SSH agent for %i

[Service]
Type=simple
ExecStart=/usr/lib/split-ssh/start-agent %i
Restart=on-failure
RestartSec=5
